# debian/postinst
#!/bin/sh
# postinst script for agave-snapshot-gossip-client

set -e

case "$1" in
    configure)
        # Install default config if it doesn't exist
        if [ ! -f /etc/gossip-client/config.toml ]; then
            cp /usr/share/doc/agave-snapshot-gossip-client/examples/config.toml /etc/gossip-client/config.toml
        fi

        # Set ownership and permissions
        chown -R gossip-client:gossip-client /etc/gossip-client
        chmod 755 /etc/gossip-client
        chmod 644 /etc/gossip-client/config.toml

        # Enable and start the service
        systemctl enable gossip-client.service
        systemctl start gossip-client.service
        ;;

    remove|deconfigure)
        # Clean up when package is removed
        systemctl stop gossip-client.service || true
        systemctl disable gossip-client.service || true
        ;;

    abort-install|abort-upgrade|abort-remove|abort-deconfigure)
        # Clean up if installation fails
        systemctl stop gossip-client.service || true
        systemctl disable gossip-client.service || true
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
