# debian/postrm
#!/bin/sh

set -e

case "$1" in
    purge)
        # Remove user, group, and config files on purge
        echo "Purging gossip-client package..."
        deluser --system --quiet gossip-client 2>/dev/null || true
        delgroup --system --quiet gossip-client 2>/dev/null || true

        # Ensure we can remove the directory
        if [ -d "/etc/gossip-client" ]; then
            echo "Removing /etc/gossip-client directory contents (preserving .json files)..."
            # Move .json files to a temporary location
            tmpdir=$(mktemp -d)
            find /etc/gossip-client -maxdepth 1 -name "*.json" -exec mv {} "$tmpdir/" \;

            # Remove the directory
            rm -rf /etc/gossip-client

            # Recreate directory and restore .json files
            mkdir -p /etc/gossip-client
            if [ -n "$(ls -A "$tmpdir/")" ]; then
                mv "$tmpdir"/* /etc/gossip-client/
                echo "Preserved .json files in /etc/gossip-client"
            fi
            rmdir "$tmpdir"

            # Verify removal of non-json files
            if [ -d "/etc/gossip-client" ]; then
                echo "Directory contents after purge:"
                ls -la /etc/gossip-client
            fi
        else
            echo "/etc/gossip-client directory not found"
        fi
        ;;
    remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        # Stop the service
        systemctl stop gossip-client.service || true
        systemctl disable gossip-client.service || true

        # Reload systemd daemon
	echo "Reloading systemd daemon"
        systemctl daemon-reload
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
