# debian/postrm
#!/bin/sh
# postrm script for agave-snapshot-gossip-client

case "$1" in
    purge)
        # Remove user, group, and config files on purge
        echo "Purging gossip-client package..."
        /usr/sbin/userdel gossip-client 2>/dev/null || true
        /usr/sbin/groupdel gossip-client 2>/dev/null || true

        # Remove non-json files from config directory
        if [ -d "/etc/gossip-client" ]; then
            echo "Removing non-json files from /etc/gossip-client..."
            find /etc/gossip-client -maxdepth 1 -type f ! -name "*.json" -delete
            echo "Directory contents after purge:"
            ls -la /etc/gossip-client
        else
            echo "/etc/gossip-client directory not found"
        fi
        ;;
    remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        # Stop the service
        systemctl stop gossip-client.service || true
        systemctl disable gossip-client.service || true

        # Reload systemd daemon
        echo "Reloading systemd daemon"
        systemctl daemon-reload
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
